# Notification/supervisor/rabbitmq-consumer.conf
[program:rabbitmq-consumer]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/artisan queue:work rabbitmq
autostart=true
autorestart=true