# Payment/supervisor/roadrunner.conf
[program:roadrunner-grpc-workers]
process_name=%(program_name)s_%(process_num)02d
command=/var/www/html/rr serve
autostart=true
autorestart=true